# scripts/dra-prod.service
[Unit]
Description=Digital Resilience Assessment — Production
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/dra/prod/server
ExecStart=/usr/bin/node -r dotenv/config server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
